# deploy/systemd/user/kdive-server.service
[Unit]
Description=KDIVE server (user)
Documentation=https://github.com/randomparity/kdive
After=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.config/kdive/kdive.env
ExecStart=%h/.local/share/kdive/.venv/bin/python -m kdive server
Restart=on-failure
RestartSec=5
# Backends (Postgres/MinIO/OIDC) are external and not systemd-ordered here; the unit
# retries until they are reachable rather than failing terminally (ADR-0114 §4).

[Install]
WantedBy=default.target
